# .docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
startretries=3
stdout_logfile=/var/log/supervisord/php-fpm.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
priority=10

[program:laravel-queue]
command=/usr/local/bin/php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=90 --queue=${REDIS_QUEUE:-default}
directory=/var/www/html
user=appuser
environment=HOME="/home/appuser",APP_ENV="local",QUEUE_CONNECTION="redis"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisord/laravel-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
redirect_stderr=true
numprocs=1
priority=20


[include]
files = /etc/supervisor/conf.d/*.conf